# deploy/duke-shibboleth-sveltekit.service
[Unit]
Description=Duke Shibboleth SvelteKit
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@@USER@@
WorkingDirectory=@@PROJECT_ROOT@@
EnvironmentFile=@@PROJECT_ROOT@@/.env
Environment=NODE_ENV=production
Environment=HOST=@@HOST@@
Environment=PORT=@@PORT@@
ExecStart=@@NODE_BIN@@ build
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=@@PROJECT_ROOT@@/certs

StandardOutput=journal
StandardError=journal
SyslogIdentifier=@@SERVICE_NAME@@

[Install]
WantedBy=multi-user.target
